# docker/root/etc/services.d/voicybot/run
#!/usr/bin/with-contenv bash

cd /app/voicybot

screen -D -m -S \
voicybot s6-setuidgid abc yarn start &

wait
