# docker/root/etc/services.d/mongodb/run
#!/usr/bin/with-contenv bash

exec /usr/bin/mongod --quiet --config /etc/mongodb.conf
